# systemd/dsmr-mqtt.service
# copy to /etc/systemd/system
# systemctl enable dsmr-mqtt.service

[Unit]
Description=P1 dsmr smartmeter

Requires=network.target
After=network.target
Wants=network-online.target systemd-networkd-wait-online.service


# The package should be installed with pip or uv
# Adapt paths and environment variables for your installation
[Service]
Type=idle
WorkingDirectory=/opt/dsmr2mqtt/
ExecStart=/usr/local/bin/dsmr2mqtt
#RemainAfterExit=true

# Environment variables for configuration
# Uncomment and set these as needed:
#Environment="MQTT_BROKER=192.168.1.1"
#Environment="MQTT_USERNAME=myuser"
#Environment="MQTT_PASSWORD=mypassword"
#Environment="SERIAL_PORT=/dev/ttyUSB0"
#Environment="DSMR_LOGLEVEL=INFO"

Restart=on-failure
RestartSec=60

[Install]
WantedBy=network-online.target
